# setup/wyliodrin-server-raspberrypi.conf
[supervisord]
[program:wyliodrind]
command=/usr/bin/wyliodrind
user=pi
autostart=true
exitcodes=255
autorestart=unexpected
environment=HOME="/wyliodrin"
priority=20
