# setup/wyliodrin-hypervisor-beagleboneblack.conf
[supervisord]
[program:wyliodrin_hypervisor]
command="/usr/bin/wyliodrin_hypervisor"
user=udooer
autostart=true
exitcodes=0
autorestart=unexpected
environment=HOME="/wyliodrin"
priority=10
